# infra/systemd/cian-scraper.timer
[Unit]
Description=Run CIAN scraper daily at 3 AM Moscow time
Documentation=https://github.com/nikitaav92-star/realestate
Requires=cian-scraper.service

[Timer]
# Run daily at 3:00 AM Moscow time (UTC+3 = 00:00 UTC in winter, 00:00 UTC+3 in summer)
# Using Europe/Moscow timezone
OnCalendar=*-*-* 00:00:00
Timezone=Europe/Moscow
# If missed, run as soon as possible
Persistent=true
# Randomize start time by up to 30 minutes to avoid thundering herd
RandomizedDelaySec=1800

[Install]
WantedBy=timers.target

